# rootfs/etc/services.d/usbip/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: VirtualHere Client Mounter
# Take down the S6 supervision tree when usbip fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "vhclientx86_64 crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "vhclientx86_64 stopped, restarting..."
